# nodepool/elements/nodepool-base/glean@.service.d/override.conf
[Service]
Environment="ARGS=--interface %I --skip-dns"
